# deploy/moraware-worker/systemd/moraware-worker.service
[Unit]
Description=eliteOS Moraware scheduled pipeline (one-shot)
Documentation=file:///opt/eliteos/eOS/docs/eliteos/moraware-cloud-worker-runbook.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=eliteos
Group=eliteos
WorkingDirectory=/opt/eliteos/eOS
EnvironmentFile=/etc/eliteos/moraware-worker.env
ExecStart=/opt/eliteos/eOS/deploy/moraware-worker/run-moraware-worker.sh
StandardOutput=append:/var/log/eliteos/moraware-nightly.log
StandardError=append:/var/log/eliteos/moraware-nightly.log

# Typical live run ~60–90 minutes; allow headroom
TimeoutStartSec=7200

[Install]
WantedBy=multi-user.target
